# resources/com.swarmgrid.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.swarmgrid.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/t/.swarmgrid/venv/bin/swarmgrid</string>
        <string>agent</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/t/clients/swarmgrid</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/t</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/tmp/swarmgrid-agent.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/swarmgrid-agent.stderr.log</string>
</dict>
</plist>
